# peq-console.service
[Unit]
Description=SpaceMaster EQ Console (port 8777)
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
WorkingDirectory=/opt/spacemaster
ExecStart=/usr/bin/python3 /opt/spacemaster/peq_toggle_nas.py
Restart=always
RestartSec=3
User=root
StandardOutput=append:/opt/spacemaster/peq-console.log
StandardError=append:/opt/spacemaster/peq-console.log

[Install]
WantedBy=multi-user.target
